# ansible/testnet-setup/rnode-monitor.service
[Service]
ExecStart=/usr/local/bin/rnode-monitor

[Install]
WantedBy=default.target
